# ops/macmini/com.myrunstreak.backfill-splits.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd agent: run one gentle splits-backfill tick every 5 minutes.

  Before loading, replace __HOME__ with your home dir (e.g. /Users/tomdrake)
  and confirm STK_BIN points at your stk binary (uv tool installs to
  ~/.local/bin/stk). Tune STK_BACKFILL_BATCH and StartInterval to taste:
  batch 5 every 300s ≈ 60 runs/hr (full history in a few days).

  Install:
    cp com.myrunstreak.backfill-splits.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.myrunstreak.backfill-splits.plist
  Stop:
    launchctl unload ~/Library/LaunchAgents/com.myrunstreak.backfill-splits.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.myrunstreak.backfill-splits</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__HOME__/myrunstreak-ops/backfill-splits.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <!-- launchd has a minimal PATH; include where stk + python live. -->
        <key>PATH</key>
        <string>__HOME__/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>STK_BIN</key>
        <string>__HOME__/.local/bin/stk</string>
        <key>STK_BACKFILL_BATCH</key>
        <string>5</string>
    </dict>

    <!-- Every 5 minutes. -->
    <key>StartInterval</key>
    <integer>300</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/myrunstreak/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/myrunstreak/launchd.err.log</string>
</dict>
</plist>
